# target/linux/ramips/mt76x8/base-files/etc/uci-defaults/90_default-wifi
#!/bin/sh
uci set wireless.radio0.disabled=0
uci set wireless.default_radio0.ssid='OpenWrt'
uci set wireless.default_radio0.encryption='psk2'
uci set wireless.default_radio0.key='Passw0rd'
uci commit wireless
wifi
